# streamserver/reboot_server.service
[Unit]
Description=Reboot Server

[Service]
Type=oneshot
ExecStart=/bin/systemctl reboot
